# provisioning/roles/docker/files/docker-tcp.socket
[Unit]
Description=Docker Socket for the API

[Socket]
ListenStream=5555
Service=docker.service
BindIPv6Only=both

[Install]
WantedBy=sockets.target
